# rootfs/etc/services.d/fail2ban/run
#!/usr/bin/execlineb -P
#

touch /tmp/empty-log.log

fail2ban-server -f
